# supervisord.conf
[supervisord]

[program:maskup]
command=node app.js
autorestart=true
stdout_logfile=server.log
